# systemd/lollipop.service
[Unit]
Description="Lollipop sticky keys service"

[Service]
Type=exec
ExecStart=/usr/bin/lollipop

[Install]
WantedBy=multi-user.target
